# templates/codex_systemd_daily.timer.template
[Unit]
Description=ZX LAB Codex Review timer — {{LABEL}}

[Timer]
OnCalendar={{ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
